# Code/process_incoming.service
# ===============================================
# Systemd Service: process_incoming.service
# -----------------------------------------------
# Zweck:
#   Führt das Skript /srv/scripts/process_incoming.sh aus,
#   das neue Uploads aus /srv/incoming_media automatisch
#   in die Zielstruktur verschiebt (z. B. /mnt/media/music, /booklets).
#
# Besonderheiten:
#   - Wird automatisch durch process_incoming.timer gestartet.
#   - Kann auch manuell per systemctl start process_incoming.service
#     ausgeführt werden.
#   - Arbeitet als One-Shot-Job (kein Dauerdienst).
#   - Schreibt alle Ausgaben in ein zentrales Logfile unter /srv/logs/.
#   - Neustart bei Fehlern (z. B. temporäre Sperren, Netzunterbrechung).
# ===============================================

[Unit]
# Beschreibung (wird in 'systemctl list-units' angezeigt)
Description=Verarbeite neue Uploads aus /srv/incoming_media

# Netzwerk muss verfügbar sein (rsync-/SSH-basierte Jobs benötigen Netzwerk)
After=network.target network-online.target
Wants=network-online.target


[Service]
# --- Allgemeine Ausführungseinstellungen ---
# "oneshot" bedeutet: das Skript läuft einmal durch und beendet sich.
Type=oneshot

# Das eigentliche Verarbeitungsskript
ExecStart=/srv/scripts/process_incoming.sh

# Benutzer und Gruppe, unter denen das Skript läuft (keine root-Rechte!)
User=medienserver
Group=medienserver

# CPU- und IO-Priorität senken, damit reguläre Benutzeraktivitäten nicht behindert werden
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Arbeitsverzeichnis setzen (optional, aber sauber)
WorkingDirectory=/srv/incoming_media

# --- Fehlerbehandlung ---
# Neustart nur bei echten Fehlern (Exit-Code ungleich 0)
Restart=on-failure
RestartSec=30s

# --- Logging ---
# Standard- und Fehlerausgabe (stdout/stderr) in Logdatei umleiten.
# Datei wird bei Bedarf automatisch erstellt.
StandardOutput=append:/srv/logs/systemd_process_incoming.log
StandardError=append:/srv/logs/systemd_process_incoming.log

# --- Umgebungsvariablen (optional, hilfreich für PATH) ---
# Falls dein Skript rsync, find, chmod o.ä. nutzt, ist PATH hier explizit gesetzt
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin"

# --- Sicherheit (optional, aber empfehlenswert) ---
# Verhindert versehentliche Schreibzugriffe außerhalb der relevanten Pfade
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true


[Install]
# Ermöglicht Aktivierung per Timer
WantedBy=multi-user.target