# Code/process_incoming.timer
# ===============================================
# Systemd Timer: process_incoming.timer
# -----------------------------------------------
# Zweck:
#   Startet regelmäßig den Service 'process_incoming.service',
#   der neue Uploads aus /srv/incoming_media verarbeitet.
#
# Vorteile:
#   - Kein Cronjob nötig (saubere systemd-Integration)
#   - Automatisches Logging & Fehlerbehandlung über systemd
#   - Flexible Steuerung: Intervall + täglicher Hauptlauf möglich
#
# Verhalten:
#   - Startet 1 Minute nach dem Bootvorgang
#   - Läuft alle 5 Minuten (um kleine Uploads zu verarbeiten)
#   - Führt zusätzlich einmal täglich um 02:00 Uhr einen
#     vollständigen Nachtlauf durch (für große Uploads)
#   - Verpasste Läufe werden nachgeholt (z. B. nach Reboot)
# ===============================================

[Unit]
# Beschreibung, erscheint in systemctl list-units und journalctl
Description=Timer für process_incoming.service – automatische Upload-Verarbeitung

[Timer]
# --- Ausführungsplan ---
# Startet 1 Minute nach dem Boot
OnBootSec=1min

# Wiederholt alle 5 Minuten (Standardintervall)
OnUnitActiveSec=5min

# Zusätzlich einmal täglich um 22:00 Uhr (Hauptlauf)
OnCalendar=*-*-* 22:00:00

# --- Stabilität ---
# Falls der Server zum geplanten Zeitpunkt ausgeschaltet war,
# wird der Job beim nächsten Start automatisch nachgeholt.
Persistent=true

# Kleine Zeitabweichung zulassen, damit nicht mehrere Timer exakt gleichzeitig starten
AccuracySec=1min

[Install]
# Aktiviert den Timer automatisch beim Booten
WantedBy=timers.target